# systemd/linux-desktop-manager.service
[Unit]
Description=Linux Desktop Manager

[Service]
Type=simple
ExecStart=/usr/local/bin/linux-desktop-manager

[Install]
WantedBy=multi-user.target
